# conf/supervisor.conf
[program:backendfail]
directory=/var/www/bf/backendfail
user=backendfail
command=/var/www/bf/env/bin/gunicorn wsgi -c /vagrant/conf/gunicorn.py
autostart=True
autorestart=True
