# conf/supervisor.conf
[program:backendfiddle]
directory=/var/www/bf/backendfiddle
user=gunicorn
command=/var/www/bf/env/bin/gunicorn wsgi -c /vagrant/conf/gunicorn.py
autostart=True
autorestart=True
